# root/etc/systemd/system/hsdepchk.service
[Unit]
After=network.target
OnFailure=mynotify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hsdepchk
User=nobody
