# home/.config/systemd/user/backup.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/.local/bin/backup"
