# home/.config/systemd/user/backup.service
[Unit]
After=network.target
OnFailure=notify@rogue:%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/.local/bin/backup"
