# home/.config/systemd/user/backup.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c %h/.local/bin/backup
